# boinc/rootfs/etc/s6-overlay/s6-rc.d/operator/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

/opt/operator/bin/operator --options /data/options.json --data /data/boinc --log-level DEBUG
