# boinc/rootfs/etc/s6-overlay/s6-rc.d/operator/up
/opt/operator/bin/operator /data/options.json
