# allow/rootfs/etc/s6-overlay/s6-rc.d/boinc/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting boinc client"
exec /usr/bin/boinc --dir /data/boinc
