# boinc/rootfs/etc/s6-overlay/s6-rc.d/boinccmd/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "TODO: Waiting boinc client start"
