# boinc/rootfs/etc/s6-overlay/s6-rc.d/boinccmd/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Waiting boinc client start"

sleep 10

bashio::log.info "Boinc client started"

# Start Account Manager configuration

cd /data/boinc || exit 1

CURRENT_ACCOUNT_MANAGER_URL=$(boinccmd --acct_mgr info | grep 'URL:' | awk '{print $2}')

if bashio::config.has_value 'account_manager_url' || bashio::config.has_value 'account_manager_username' ||  bashio::config.has_value 'account_manager_password'; then

    bashio::log.info "Start account manager configuration"

    bashio::config.require 'account_manager_url'
    ACCOUNT_MANAGER_URL="$(bashio::config 'account_manager_url')"

    bashio::config.require.username 'account_manager_username'
    ACCOUNT_MANAGER_USERNAME="$(bashio::config 'account_manager_username')"

    bashio::config.require.password 'account_manager_password'
    ACCOUNT_MANAGER_PASSWORD="$(bashio::config 'account_manager_password')"

    if [ -z "$CURRENT_ACCOUNT_MANAGER_URL" ]; then
        bashio::log.info "Attaching to $ACCOUNT_MANAGER_URL"
        boinccmd --acct_mgr attach "$ACCOUNT_MANAGER_URL" "$ACCOUNT_MANAGER_USERNAME" "$ACCOUNT_MANAGER_PASSWORD"
    else
        bashio::log.info "Attached to $CURRENT_ACCOUNT_MANAGER_URL"
        if [ "${ACCOUNT_MANAGER_URL%/}" = "${CURRENT_ACCOUNT_MANAGER_URL%/}" ]; then
            bashio::log.info "Syncronizing from $CURRENT_ACCOUNT_MANAGER_URL"
            boinccmd --acct_mgr sync
        else
            bashio::log.info "Detachhing from $CURRENT_ACCOUNT_MANAGER_URL"
            boinccmd --acct_mgr detach
            sleep 10
            bashio::log.info "Attaching to $ACCOUNT_MANAGER_URL"
            boinccmd --acct_mgr attach "$ACCOUNT_MANAGER_URL" "$ACCOUNT_MANAGER_USERNAME" "$ACCOUNT_MANAGER_PASSWORD"
        fi
    fi

    bashio::log.info "Account manager configured"
else
    if [ -n "$CURRENT_ACCOUNT_MANAGER_URL" ]; then
        bashio::log.info "Deteching configured account manager"
        bashio::log.info "Detachhing from $CURRENT_ACCOUNT_MANAGER_URL"
        boinccmd --acct_mgr detach
    fi
fi

# End Account Manager configuration
