# boinc/rootfs/etc/s6-overlay/s6-rc.d/config/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Configuring Boinc Addon"

if ! bashio::fs.directory_exists "/data/boinc"; then
    bashio::log.info "Creating boinc folders"
    mkdir -p /data
    mkdir -p /data/boinc
    mkdir -p /data/boinc/slots
fi

HOSTNAME=$(bashio::info.hostname)

if [ "$HOSTNAME" == 'null' ]; then
    HOSTNAME="homeassistant"
fi

bashio::log.info "Hostname $HOSTNAME configured"

echo "{\"hostname\": \"$HOSTNAME\"}" | tempio \
    -template /usr/share/tempio/cc_config.xml \
    -out /data/boinc/cc_config.xml

bashio::log.info "Boinc Addon configured"
