# boinc/rootfs/etc/services.d/boinc/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting boinc home assistant addon"

if ! bashio::fs.directory_exists "/data/boinc"; then
    bashio::log.info "Creating boinc folders"
    mkdir -p /data
    mkdir -p /data/boinc
    mkdir -p /data/boinc/slots
fi

HOSTNAME=$(bashio::info.hostname)

if [ "$HOSTNAME" == 'null' ]; then
    HOSTNAME="homeassistant"
fi

echo "{\"hostname\": \"$HOSTNAME\"}" | tempio \
    -template /usr/share/tempio/cc_config.xml \
    -out /data/boinc/cc_config.xml

PROJECT="$(bashio::config 'project')"
KEY="$(bashio::config 'key')"

bashio::log.info "Starting boinc client"

exec /usr/bin/boinc --dir /data/boinc --no_gui_rpc --attach_project "$PROJECT" "$KEY" --update_prefs "$PROJECT"