# boinc/rootfs/etc/services.d/boinc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting boinc home assistant addon"

bashio::log.info "Creating boinc folders"
mkdir -p /data
mkdir -p /data/boinc
mkdir -p /data/boinc/slots

bashio::log.info "Starting boinc client"

PROJECT="$(bashio::config 'project')"
KEY="$(bashio::config 'key')"


exec /usr/bin/boinc --dir /data/boinc --no_gui_rpc --attach_project "$PROJECT" "$KEY" --update_prefs "$PROJECT"